# arm/etc/rc.local
#!/bin/sh -e
echo "Regenerating host keys ..." >> /dev/kmsg

# Drain entropy pool.
dd if=/dev/urandom of=/dev/null bs=1024 count=10 2>/dev/null

# Regenerate the keys.
rm -f /etc/ssh/ssh_host_*
dpkg-reconfigure openssh-server

echo "#!/bin/sh -e
exit 0
" > /etc/rc.local

exit 0
